# _obs_audit_20251030_214213/_systemd/trading-metrics.service
[Unit]
Description=Trading Prometheus Metrics Sidecar (:9110)
After=network-online.target redis.service
Wants=network-online.target

[Service]
User=trader
Group=trader
EnvironmentFile=/srv/trading/.env.trading
WorkingDirectory=/srv/trading
ExecStart=/srv/trading/.venv/bin/python /srv/trading/tools/metrics_sidecar.py
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
